# free-ask-internet/supervisor.d/chatgpt-next-web.ini
[program:chatgpt-next-web]
directory=/app/chatgpt-next-web
command=node server.js
environment=BASE_URL="http://backend:8000",OPENAI_API_KEY="EMPTY",CUSTOM_MODELS="-all,+gpt-3.5-turbo"
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true