# free-ask-internet/supervisor.d/proxy.ini
[program:proxy]
directory=/app/proxy
command=/app/proxy/proxy
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true