# ai/next-chat-cf/supervisor.d/workers-ai-to-openai.ini
[program:workers-ai-to-openai]
directory=/app
command=/app/workers-ai-to-openai/workers-ai-to-openai
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true